# ops/systemd/flywheel-sweep.timer
[Unit]
Description=Nightly food-cache flywheel sweep

[Timer]
OnCalendar=*-*-* 04:30
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
